# v1/fleet/required/update-os.timer
[Unit]
Description=Runs CoreOS update check every 5m

[Timer]
OnBootSec=10min
OnUnitActiveSec=5min

[X-Fleet]
Global=true
